# ansible/roles/medley/templates/medley.service
[Service]
WorkingDirectory=/srv/medley
ExecStart=/srv/medley/venv/bin/python medley.py
Type=forking
PIDFile=/var/run/medley.pid
