# ansible/roles/application/files/medley-headlines.timer
[Unit]
Description=Timer for Medley headlines precache

[Timer]
# Daily at 7:00 AM
OnCalendar=07:00

[Install]
WantedBy=timers.target
