# ansible/roles/application/files/medley-geodb.timer
[Unit]
Description=Timer for Medley geodb app

[Timer]
# Second Monday of the month at 1 AM
OnCalendar=Monday *-*-8..14 01:00:00

# Catch up if schedule was missed
Persistent=true

[Install]
WantedBy=timers.target
