# ansible/roles/application/files/medley-fs-maintenance.timer
[Unit]
Description=Timer for Medley filesystem maintenance

[Timer]
# Daily at 5:03 AM
OnCalendar=05:03

[Install]
WantedBy=timers.target
