# ansible/roles/application/templates/medley-failure@.service.j2
[Unit]
Description=Medley Failure Notifier

[Service]
Type=oneshot
ExecStart={{ project_root }}/venv/bin/pthon {{ project_root }}/scripts/failmail.py %i {{ admin_email }}

# Local Variables:
# mode: conf
# End:
